# beam/BEAM_supervisor_systemd_package/beam-daily-sync.timer
[Unit]
Description=Run BEAM daily data sync once per day

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=beam-daily-sync.service

[Install]
WantedBy=timers.target
