# beam/BEAM_supervisor_systemd_package/beam-health-check.timer
[Unit]
Description=Run BEAM health check every 10 minutes

[Timer]
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true
Unit=beam-health-check.service

[Install]
WantedBy=timers.target
